# etc/network/if-up.d/noip-update
#!/bin/sh

#=====================================================================
# CONFIGURATION VARIABLES
#
# Change these to whatever you registered at noip.com, put this script
# in '/etc/network/if-up.d/noip-update', and give it rwx permissions
# for root only (e.g. chmod 0700 noip-update) because it contains your
# noip.com password in clear text.
#
#---------------------------------------------------------------------
NOIP_USER='my_no-ip_username'
NOIP_PASS='my_no-ip_password'
NOIP_HOST='my.no-ip.host.name'
#=====================================================================


#=====================================================================
# LOG ROTATION ('logrotate' configuration, optional)
#---------------------------------------------------------------------
# The following 'logrotate' configuration (without the leading hashes,
# of course) may be added to '/etc/logrotate.d/noip-update' so the log
# file generated by this script gets rotated for a few weeks and
# eventually deleted so that the log file doesn't just grow
# indefinitely:
#
# /var/log/noip-update.log {
#        weekly
#        missingok
#        rotate 5
#        compress
#        create 640 root adm
#        notifempty
# }
#
#=====================================================================


#=======================================================
# Probably no need to change anything beyond this point.
#-------------------------------------------------------
NOIP_STUB='dynupdate.no-ip.com/nic/update'
SELF="$(basename $0)"

log () {
    echo "[$(date +'%Y-%m-%d %H:%M:%S')]: $*" >> "${LOGFILE}"
}

url () {
    echo "https://${NOIP_USER}:${NOIP_PASS}@${NOIP_STUB}?hostname=${NOIP_HOST}&myip=${DHCP4_IP_ADDRESS}"
}

if [ -z "${DHCP4_IP_ADDRESS}" ]; then
    log 'No IP address supplied... exiting'
    exit 1
fi

if ! which curl > /dev/null; then
    log 'No "curl" on PATH... exiting'
    exit 2
fi

response=$(curl "$(url)" 2>/dev/null)
log "${DHCP4_IP_ADDRESS} RESPONSE: ${response}"
